# ryzenadj-linux/ryzenboost.service
[Unit]
Description=Ryzenboost

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/bin/ryzenboost

[Install]
WantedBy=multi-user.target
