# ryzenadj-linux/ryzenboost.service
[Unit]
Description=Ryzenboost

[Service]
Type=simple
ExecStart=/usr/bin/ryzenboost

[Install]
WantedBy=multi-user.target
